# systemd/user/bspwm.target
[Unit]
Description=BSPWM start services that need graphical-session.target
Requires=wm.target
BindsTo=wm.target
