# systemd/user/docker-2w-update.timer
[Unit]
Description=All dockers container updater timer

[Timer]
OnUnitActiveSec=2w
Persistent=true

[Install]
WantedBy=timers.target
