# systemd/user/searxng-update.timer
[Unit]
Description=Searxng container updater timer

[Timer]
OnActiveSec=2w

[Install]
WantedBy=timers.target
